# sunxi/pinephone-boot-eg25/autobuild/overrides/usr/lib/systemd/system/boot-eg25.service
[Unit]
Description=Enable EG25 modem on PinePhone V1.0 prototype

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/boot_eg25

[Install]
WantedBy=network.target
