# omap3-n900/enable-n900-modem/autobuild/overrides/usr/lib/systemd/system/n900-modem.service
[Unit]
Description=Enable the modem of N900
Before=ofono.service

[Service]
Type=oneshot
ExecStart=/usr/bin/enable-n900-modem

[Install]
WantedBy=multi-user.target
