# arm-qcom/firmware-xiaomi-wt88047/autobuild/overrides/usr/lib/systemd/system/mnt-persist.mount
[Unit]
Description=Mount persistent storage
Requires=systemd-fsck@dev-disk-by\x2dpartlabel-persist.service
After=systemd-fsck@dev-disk-by\x2dpartlabel-persist.service
After=blockdev@dev-disk-by\x2dpartlabel-persist.target

[Mount]
Where=/mnt/persist
What=/dev/disk/by-partlabel/persist
Type=ext4
Options=ro
