# sunxi-old/sun8iw5p1-workaround/autobuild/overrides/usr/lib/systemd/system/sun8iw5p1-workaround.service
[Unit]
Description=SUN8IW5P1 Workaround

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sun8iw5p1-workaround start

[Install]
WantedBy=multi-user.target
